# .docker/postgrest/crontab
# Run update_jwk.sh at midnight
0 0 */1 * * /app/update_jwk.sh >> /var/log/cron.log 2>&1
